# src/dbus-org.nemomobile.provisioning.service
[Unit]
Description=Provisioning service

[Service]
Type=dbus
BusName=org.nemomobile.provisioning
ExecStart=/usr/libexec/provisioning-service
